# ttyd/files/ttyd.service.jinja
[Unit]
Description=ttyd
After=network.target

[Service]
User=root
ExecStart=/usr/local/bin/ttyd -p {{ port }} login
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
